# scripts/cine.service
[Unit]
Description=Cine Node.js App
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/cine
ExecStart=/usr/bin/node /opt/cine/app.js
Restart=on-failure
RestartSec=10
Environment="PORT=3000"
EnvironmentFile=-/etc/cine/cine.env

[Install]
WantedBy=multi-user.target
